# image/rootfs-overlay/etc/systemd/system/pivideo.service
[Unit]
Description=PiVideo GPIO Button Daemon
After=multi-user.target

[Service]
ExecStartPre=/bin/dmesg -n 1
ExecStart=/usr/local/bin/pivideo-daemon
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info
Environment=PIVIDEO_CONFIG=/opt/pivideo/config.json
Environment=PIVIDEO_VIDEO_DIR=/opt/pivideo/videos
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
